# app/config/config.json
{
  "mongodb": {
    "db_name": "siem_project",
    "uri": "mongodb://localhost:27017"
  },
  "redis": {
    "url": "redis://localhost:6379",
    "key_prefix": "warsoc:stateful:",
    "default_ttl": 3600
  },
  "threat_intelligence": {
    "enabled": true,
    "ips": ["192.168.1.100", "10.0.0.200"],
    "files": ["data/blacklist_ip.txt"]
  },
  "whitelist": {
    "enabled": true,
    "service_accounts": ["svc_backup", "nagios", "backup", "root", "daemon"],
    "ips": ["127.0.0.1", "::1", "localhost"]
  },
  "geoip": {
    "enabled": true,
    "mock_mode": true
  },
  "detection": {
    "brute_force_threshold": 3,
    "failed_login_patterns": ["failed password", "authentication failure", "login failed", "access denied"],
    "rules": {
      "SQL_INJECTION": {"regex": "(%27|')\\s*(union|select|insert|update|delete|drop|or\\s+1\\s*=\\s*1|--|#)", "sev": "HIGH", "mitre": "T1190"},
      "XSS_ATTACK": {"regex": "(?i)(<script|javascript:|onerror=|onload=)", "sev": "HIGH", "mitre": "T1190"},
      "COMMAND_INJECTION": {"regex": "(?i)(;\\s*whoami|\\|\\s*nc|\\$\\s*\\(|`whoami`)", "sev": "CRITICAL", "mitre": "T1059"},
      "PATH_TRAVERSAL": {"regex": "(?i)(\\.\\.\\/|\\.\\.\\\\|etc/passwd|win.ini)", "sev": "HIGH", "mitre": "T1083"},
      "POWERSHELL_OBFUSCATION": {"regex": "powershell.*(\\-e(nc)?(odedcommand)?|\\-w(indowstyle)?\\s+h(idden)?)\\s+[A-Za-z0-9+/=]{20,}", "sev": "CRITICAL", "mitre": "T1027"},
      "PRIVILEGE_ESCALATION": {"regex": "(sudo\\s+su\\s+\\-|nopasswd:\\s+all|net\\s+localgroup\\s+administrators\\s+.*\\s+/add|passwd\\s+\\-d)", "sev": "HIGH", "mitre": "T1068"},
      "LATERAL_MOVEMENT": {"regex": "(net\\s+use\\s+\\\\\\\\[a-z0-9.]+\\\\c\\$|wmic\\s+/node:.*process\\s+call\\s+create|winrm\\s+invoke|enter-pssession)", "sev": "HIGH", "mitre": "T1021"},
      "LOG_EVASION": {"regex": "(wevtutil\\s+cl\\s+|history\\s+\\-c|rm\\s+.*\\.log|shred\\s+|clear-eventlog|set-mppreference\\s+-disablerealtime)", "sev": "HIGH", "mitre": "T1070"},
      "REVERSE_SHELL": {"regex": "(bash\\s+\\-i\\s+>\\&\\s+/dev/tcp/|python.*\\-c.*socket|perl\\s+\\-e.*socket|nc\\s+-e\\s+/bin/sh)", "sev": "CRITICAL", "mitre": "T1059.004"},
      "RECON_COMMANDS": {"regex": "(?i)(whoami|ipconfig|netstat|enum4linux|nmap|arp -a)", "sev": "MEDIUM", "mitre": "T1592"},
      "PERSISTENCE": {"regex": "(?i)(crontab -e|authorized_keys|registry run keys|schtasks /create)", "sev": "HIGH", "mitre": "T1053"},
      "DATA_EXFILTRATION": {"regex": "(?i)(wget|curl|tar -czf|zip -r|nc -w|sftp|scp)", "sev": "CRITICAL", "mitre": "T1041"},
      "STAGING": {"regex": "(7z\\s+a|zip\\s+\\-r|tar\\s+\\-czf)\\s+.*(/tmp/|/windows/temp/|/appdata/|/var/tmp/)", "sev": "MEDIUM", "mitre": "T1074"},
      "BRUTE_FORCE_PATTERN": {"regex": "(?i)(failed password|authentication failure|login failed|access denied)", "sev": "MEDIUM", "mitre": "T1110"},
      "XXE_INJECTION": {"regex": "(?i)(<!ENTITY|SYSTEM|file:///)", "sev": "CRITICAL", "mitre": "T1190"},
      "MALWARE_EXECUTION": {"regex": "(?i)(mimikatz|hashcat|psexec|vssadmin|net user)", "sev": "CRITICAL", "mitre": "T1059"}
    }
  },
  "stateful_detection_rules": {
    "auth_identity": {
      "high_velocity_brute_force": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1110.001",
        "severity": "HIGH",
        "description": "High-velocity brute force attack detected",
        "group_by": "source_ip",
        "event_filter": "failed_login"
      },
      "low_slow_brute_force": {
        "enabled": true,
        "threshold": 20,
        "window_seconds": 3600,
        "mitre_id": "T1110.001",
        "severity": "MEDIUM",
        "description": "Low-and-slow brute force attack detected",
        "group_by": "source_ip",
        "event_filter": "failed_login"
      },
      "password_spraying": {
        "enabled": true,
        "threshold_users": 10,
        "threshold_attempts": 1,
        "window_seconds": 300,
        "mitre_id": "T1110.003",
        "severity": "HIGH",
        "description": "Password spraying attack detected",
        "group_by": "source_ip",
        "unique_field": "username",
        "event_filter": "failed_login"
      },
      "impossible_travel": {
        "enabled": true,
        "max_travel_time_hours": 2,
        "min_distance_km": 500,
        "mitre_id": "T1078",
        "severity": "CRITICAL",
        "description": "Impossible travel detected",
        "event_filter": "successful_login"
      },
      "new_location_access": {
        "enabled": true,
        "mitre_id": "T1078",
        "severity": "MEDIUM",
        "description": "Login from new geographic location",
        "event_filter": "successful_login"
      },
      "concurrent_sessions": {
        "enabled": true,
        "threshold": 3,
        "window_seconds": 300,
        "mitre_id": "T1078",
        "severity": "MEDIUM",
        "description": "Concurrent sessions from multiple IPs",
        "group_by": "username",
        "event_filter": "successful_login"
      },
      "after_hours_activity": {
        "enabled": true,
        "start_hour": 2,
        "end_hour": 5,
        "mitre_id": "T1078",
        "severity": "MEDIUM",
        "description": "After-hours suspicious activity",
        "event_filter": "successful_login"
      },
      "account_storm": {
        "enabled": true,
        "threshold": 5,
        "window_seconds": 60,
        "mitre_id": "T1136",
        "severity": "HIGH",
        "description": "Mass account creation detected",
        "event_filter": "account_created"
      },
      "privilege_escalation_spike": {
        "enabled": true,
        "threshold": 5,
        "window_seconds": 60,
        "mitre_id": "T1548",
        "severity": "HIGH",
        "description": "Privilege escalation spike detected",
        "event_filter": "privilege_escalation",
        "keywords": ["sudo", "runas", "su", "elevation"]
      },
      "dormant_account_activation": {
        "enabled": true,
        "dormant_days": 30,
        "mitre_id": "T1078.002",
        "severity": "MEDIUM",
        "description": "Dormant account reactivated",
        "event_filter": "successful_login"
      }
    },
    "filesystem_ransomware": {
      "mass_file_modification": {
        "enabled": true,
        "threshold": 50,
        "window_seconds": 60,
        "mitre_id": "T1486",
        "severity": "CRITICAL",
        "description": "Mass file modification - Ransomware indicator",
        "group_by": "username",
        "event_filter": "file_modify"
      },
      "mass_file_deletion": {
        "enabled": true,
        "threshold": 50,
        "window_seconds": 60,
        "mitre_id": "T1485",
        "severity": "CRITICAL",
        "description": "Mass file deletion detected",
        "group_by": "username",
        "event_filter": "file_delete"
      },
      "sensitive_file_touch": {
        "enabled": true,
        "threshold": 5,
        "window_seconds": 60,
        "mitre_id": "T1005",
        "severity": "HIGH",
        "description": "Multiple sensitive file accesses",
        "group_by": "username",
        "event_filter": "file_access",
        "sensitive_keywords": ["confidential", "finance", "payroll", "secret", "private"]
      },
      "ransomware_extensions": {
        "enabled": true,
        "threshold": 1,
        "window_seconds": 60,
        "mitre_id": "T1486",
        "severity": "CRITICAL",
        "description": "Ransomware file extension detected",
        "group_by": "username",
        "event_filter": "file_write",
        "ransomware_extensions": [".crypt", ".locked", ".encrypted", ".locky", ".wannacry", ".cerber", ".cryptolocker"]
      },
      "data_exfiltration_volume": {
        "enabled": true,
        "threshold_bytes": 104857600,
        "window_seconds": 300,
        "mitre_id": "T1041",
        "severity": "HIGH",
        "description": "Large data upload detected",
        "group_by": "username",
        "event_filter": "data_upload"
      },
      "log_clearing_sequence": {
        "enabled": true,
        "window_seconds": 60,
        "mitre_id": "T1070",
        "severity": "CRITICAL",
        "description": "Log clearing sequence detected",
        "group_by": "source_ip",
        "event_sequence": ["stop_logging", "clear_logs"]
      }
    },
    "network_lateral": {
      "vertical_port_scan": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1046",
        "severity": "HIGH",
        "description": "Vertical port scan detected",
        "group_by": "source_ip",
        "unique_field": "destination_port",
        "event_filter": "network_connection"
      },
      "horizontal_port_scan": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1046",
        "severity": "HIGH",
        "description": "Horizontal port scan detected",
        "group_by": "source_ip",
        "unique_field": "destination_ip",
        "event_filter": "network_connection"
      },
      "beaconing_c2": {
        "enabled": true,
        "threshold": 10,
        "variance_seconds": 1,
        "mitre_id": "T1071",
        "severity": "CRITICAL",
        "description": "C2 beaconing pattern detected",
        "group_by": "source_ip",
        "event_filter": "network_connection"
      },
      "long_duration_connection": {
        "enabled": true,
        "threshold_hours": 24,
        "mitre_id": "T1572",
        "severity": "MEDIUM",
        "description": "Long-duration connection - Potential tunnel",
        "event_filter": "network_connection"
      },
      "rare_port_usage": {
        "enabled": true,
        "mitre_id": "T1571",
        "severity": "MEDIUM",
        "description": "Rare port usage detected",
        "event_filter": "network_connection",
        "suspicious_ports": [6667, 6697, 31337, 12345, 54321, 1337]
      },
      "dns_tunneling": {
        "enabled": true,
        "threshold": 50,
        "window_seconds": 60,
        "subdomain_min_length": 50,
        "mitre_id": "T1071.004",        "severity": "HIGH",
        "description": "DNS tunneling detected",
        "group_by": "source_ip",
        "event_filter": "dns_query"
      },
      "rdp_brute_force": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1021.001",
        "severity": "HIGH",
        "description": "RDP brute force detected",
        "group_by": "source_ip",
        "event_filter": "rdp_failed",
        "target_port": 3389
      },
      "smb_storm": {
        "enabled": true,
        "threshold": 100,
        "window_seconds": 60,
        "mitre_id": "T1021.002",
        "severity": "HIGH",
        "description": "SMB access storm detected",
        "group_by": "source_ip",
        "event_filter": "smb_access",
        "target_port": 445
      }
    },
    "web_application": {
      "high_vol_sql_injection": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1190",
        "severity": "HIGH",
        "description": "SQL injection flood detected",
        "group_by": "source_ip",
        "event_filter": "http_request",
        "sql_patterns": ["'", "union", "select", "insert", "drop", "delete", "--", "/*"]
      },
      "directory_brute_force": {
        "enabled": true,
        "threshold": 20,
        "window_seconds": 60,
        "mitre_id": "T1083",
        "severity": "MEDIUM",
        "description": "Directory brute force detected",
        "group_by": "source_ip",
        "event_filter": "http_404"
      },
      "waf_evasion_flood": {
        "enabled": true,
        "threshold": 50,
        "window_seconds": 60,
        "max_header_size": 8192,
        "max_payload_size": 1048576,
        "mitre_id": "T1190",
        "severity": "HIGH",
        "description": "WAF evasion flood detected",
        "group_by": "source_ip",
        "event_filter": "http_request"
      },
      "xss_flooding": {
        "enabled": true,
        "threshold": 10,
        "window_seconds": 60,
        "mitre_id": "T1059.007",
        "severity": "MEDIUM",
        "description": "XSS flood detected",
        "group_by": "source_ip",
        "event_filter": "http_request",
        "xss_patterns": ["<script", "javascript:", "onerror=", "onload="]
      },
      "bot_scraping": {
        "enabled": true,
        "threshold": 100,
        "window_seconds": 60,
        "mitre_id": "T1213",
        "severity": "LOW",
        "description": "Bot scraping detected",
        "group_by": "source_ip",
        "event_filter": "http_request",
        "target_path_prefix": "/product"
      },
      "error_storm": {
        "enabled": true,
        "threshold": 50,
        "window_seconds": 60,
        "mitre_id": "T1499",
        "severity": "MEDIUM",
        "description": "Application error storm detected",
        "group_by": "source_ip",
        "event_filter": "http_500"
      }
    }
  }
}